# crontab
40 * * * * bash /home/nandhu/Documents/notes/sync.sh
